# scripts/coco/podvm/luks-scratch/etc/systemd/system/kata-agent.service.d/10-override.conf
[Service]
ExecStartPre=sh -c '[ -b /dev/mapper/scratch ] && mount /dev/mapper/scratch /run/kata-containers'
