# konflux/podvm-root/etc/systemd/system/kata-agent.service.d/10-override.conf
[Service]
ExecStartPre=sh -c '[ -b /dev/mapper/scratch ] && mount /dev/mapper/scratch /kata-containers'
